# fm.service
[Unit]
Description=FM Exiter

[Service]
ExecStart=fm
After=radio.service,rds_song.service
Wants=radio.service

[Install]
WantedBy=multi-user.target
