# fm.service
[Unit]
Description=FM Exiter

[Service]
ExecStart=/srv/media/radio/fm.py

[Install]
WantedBy=multi-user.target
