# rds_song.service
[Unit]
Description=RDS update script
Wants=fm.service
[Service]
ExecStart=rds_song
Before=fm.service
[Install]
WantedBy=multi-user.target
